# i.conf
[program:i]
command=python3 /home/curator2/i/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/i/test.err.log
stdout_logfile=/home/curator2/i/test.out.log